# systemd/kuro-exporter.timer
[Unit]
Description=Run Kuro data exporter every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
